# systemd/.config/systemd/user/arr-stall-requeue.timer
[Unit]
Description=Roda arr-stall-requeue a cada 30 minutos

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
