# systemd/.config/systemd/user/xwayland-x1-symlink.service
[Unit]
Description=Keep XWayland X1 socket symlink alive (Hyprland creates X1_ instead of X1)
Documentation=https://github.com/hyprwm/Hyprland
After=graphical-session.target

[Service]
Type=simple
# Hyprland's XWayland binds the socket as /tmp/.X11-unix/X1_ ; Steam and other
# X clients look for X1. The symlink vanishes when Steam restarts XWayland
# without restarting Hyprland. Poll (cheap) and re-create it whenever it is
# missing or pointing at the wrong target.
ExecStart=/bin/bash -c 'while :; do \
  if [ -S /tmp/.X11-unix/X1_ ] && [ "$(readlink /tmp/.X11-unix/X1 2>/dev/null)" != "/tmp/.X11-unix/X1_" ]; then \
    ln -sf /tmp/.X11-unix/X1_ /tmp/.X11-unix/X1; \
  fi; \
  sleep 3; \
done'
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
